# cron/src/test/resources/main.crontab
*/10 * * * * ? echo "Test"
#*/5 * * * * ? echo "Test2"